# coc/docker/systemd/system/multi-user.target.wants/coc-ws-market.service
[Unit]
Description=Websocket for market
Requires=coc-memcached.service

[Service]
Type=simple
Restart=always
RestartSec=1
WorkingDirectory=/var/www/html/coc
ExecStart=bash start_market.sh

[Install]
WantedBy=multi-user.target
